# etc/cron/crontab

#
# Ensure that our search indexes are current.
#
12 0 * * * cd ~/current && ./bin/build-search-index  >/dev/null
